# crontab
*/1       *       *       *       *       sh /usr/local/bin/backup.sh